# root/etc/s6-overlay/s6-rc.d/svc-cdp-forward/run
#!/command/with-contenv bash
exec socat TCP-LISTEN:9222,reuseaddr,fork TCP:127.0.0.1:9666
